# supervisord.conf
; ============================================================
; supervisord.conf — gère les 3 agents FastAPI
; ============================================================

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ── 1. Water Agent (port 8001) ────────────────────────────
[program:water_agent]
command=python agents/water_agent_server.py
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/var/log/supervisor/water_agent.log
stderr_logfile=/var/log/supervisor/water_agent_err.log
environment=PYTHONIOENCODING="utf-8",PYTHONUTF8="1"

; ── 2. Manure Agent (port 8002) ───────────────────────────
[program:manure_agent]
command=python agents/manure_agent_server.py
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/var/log/supervisor/manure_agent.log
stderr_logfile=/var/log/supervisor/manure_agent_err.log
environment=PYTHONIOENCODING="utf-8",PYTHONUTF8="1"

; ── 3. Master + Dashboard (port 8000) ─────────────────────
[program:master]
command=python main_a2a.py
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=10
; Attendre que les sous-agents soient prêts
depends_on=water_agent,manure_agent
stdout_logfile=/var/log/supervisor/master.log
stderr_logfile=/var/log/supervisor/master_err.log
environment=PYTHONIOENCODING="utf-8",PYTHONUTF8="1",MASTER_PORT="8000"